# ssss.conf
[program:ssss]
command=python3 /home/curator2/ssss/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ssss/test.err.log
stdout_logfile=/home/curator2/ssss/test.out.log